# pondtv/config/pondtv.service
[Unit]
Description=PondTV - Zen Media Player
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pondtv
ExecStart=/usr/bin/python3 /home/pi/pondtv/run_pondtv.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target 