# packaging/pondtv.service
# PondTV systemd unit. Install with:
#
#   sudo cp packaging/pondtv.service /etc/systemd/system/pondtv.service
#   sudo systemctl enable --now pondtv
#
# It owns tty1 so mpv can grab DRM/KMS, and restarts on any crash (the manager
# also reconnects to mpv on its own, but a hard process death is systemd's job).
[Unit]
Description=PondTV dumb-TV channel
After=systemd-user-sessions.service
# We take over tty1 from the login getty.
Conflicts=getty@tty1.service

[Service]
Type=simple
# Needs root for mount/umount, raw evdev, and DRM/KMS on the console.
User=root
# Kiosk: a guest pressing Esc/Q must not drop to a black console.
Environment=PONDTV_ALLOW_QUIT=0
# Own the console so mpv can use /dev/dri/card* via DRM/KMS.
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
# Run from the directory that contains the `pondtv` package so `-m pondtv`
# resolves it. Adjust if you `pip install` the package instead (then CWD is
# irrelevant and you can drop WorkingDirectory).
WorkingDirectory=/home/kdklv
ExecStart=/usr/bin/python3 -m pondtv
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
