# dist/payload/pondtv.service
[Unit]
Description=PondTV Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/pondtv/run.py
WorkingDirectory=/opt/pondtv
StandardOutput=journal
StandardError=journal
Restart=always
User=pi

[Install]
WantedBy=multi-user.target 